# ipv6-tunnel/etc/cron.d/root
*/1 * * * * /bin/sh /opt/he-ipv6/update_tunnel.sh
